# deploy/discord-stat-bot-save.service
[Unit]
Description=forces stat-bot to save

[Service]
Type=oneshot
ExecStart=/usr/bin/docker kill --signal="SIGINT" stat-bot
